# ota-usb-flash/usb-flash.service
[Unit]
Description=test python script

[Service]
Type=simple
ExecStart=/usr/local/bin/usb-reflash.py

[Install]
WantedBy=multi-user.target
